# scripts/miranda.service
[Unit]
Description=Miranda - Telegram bot for remote Claude orchestration
After=network.target

[Service]
Type=simple
# IMPORTANT: Must run as non-root user. Claude's --dangerously-skip-permissions
# is blocked when running as root for security reasons.
User=%i
Group=%i
WorkingDirectory=/home/%i/miranda
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PATH=/home/%i/.cargo/bin:/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/%i/.config/miranda/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=miranda

[Install]
WantedBy=multi-user.target
